# services/plists/com.r.sync-ueli-launchers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.r.sync-ueli-launchers</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/r/.config/bin/sync-ueli-launchers</string>
    </array>

    <!-- Run every 60 seconds to keep launchers in sync with Chrome profiles -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <!-- Timeout after 30 seconds to prevent hangs -->
    <key>ExitTimeOut</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/r/.config/services/logs/sync-ueli-launchers.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/r/.config/services/logs/sync-ueli-launchers.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
