# services/plists/com.r.thegrid-dev.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.r.thegrid-dev</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/r/repos/theGrid/grid-server/.build/debug/GridServer.app/Contents/MacOS/grid-server</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/r/repos/theGrid</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/r/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/r/.config/services/logs/thegrid-dev.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/r/.config/services/logs/thegrid-dev.log</string>
</dict>
</plist>
